# aarchup/.config/systemd/system/aarchup.service
[Unit]
Description=Run aarchup

[Service]
Type=simple
ExecStart=/usr/bin/aarchup -t 5 -a "yay -Qu"
