# aarchup/.config/systemd/system/aarchup.timer
[Unit]
Description=Run aarchup in a given frequency 

[Timer]
OnBootSec=5m
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
